# roles/ceph/monitor/templates/ceph-monitor.service.j2
#
#   Author: Rohith
#   Date: 2015-06-11 16:17:51 +0100 (Thu, 11 Jun 2015)
#
#  vim:ts=2:sw=2:et
#

[Unit]
Description=Ceph Monitor
After=docker.service

[Service]
EnvironmentFile=/etc/environment
Restart=always
RestartSec=10s
TimeoutStartSec=0
TimeoutStopSec=0

ExecStartPre=/usr/bin/docker pull {{ ceph_monitor_image }}:{{ ceph_monitor_image_tag }}
ExecStartPre=-/usr/bin/docker kill ceph-monitor
ExecStartPre=-/bin/mkdir -p /var/lib/ceph
ExecStart=/usr/bin/bash -c "/usr/bin/docker run \
  --rm --net=host --name=ceph-monitor \
  -e HOST=${COREOS_PRIVATE_IPV4} \
  -e ETCD_PORT=4001 \
  -e ETCD=${COREOS_PRIVATE_IPV4} \
  -e ETCD_PATH=/deis/store \
  -e HOSTNAME=${HOSTNAME} \
  -v /etc/ceph:/etc/ceph \
  -v /var/lib/ceph:/var/lib/ceph \
   {{ ceph_monitor_image }}:{{ ceph_monitor_image_tag }}"

[Install]
WantedBy=multi-user.target

[X-Fleet]
MachineMetadata=ceph_monitor=true
Global=true
