# roles/kubernetes/common/templates/profile.d/kubernetes.sh.j2
#!/bin/bash
#
#   Author: Rohith
#   Date: 2015-05-03 17:33:05 +0100 (Sun, 03 May 2015)
#
#  vim:ts=2:sw=2:et
#

KUBERNETES="{{ kubernetes_dir }}"
export PATH=${PATH}:${KUBERNETES}/bin

